# .config/systemd/user/xinit.target
[Unit]
Description=xinit session running
BindsTo=graphical-session.target

[Install]
after=graphical-session.target
